# azure-vm/colpali.service
[Unit]
Description=ColPali Document Embedding API
Documentation=https://github.com/your-org/colpali-docker
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/colpali

# Start the services
ExecStart=/usr/local/bin/docker-compose up -d

# Stop the services
ExecStop=/usr/local/bin/docker-compose down

# Restart policy
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=colpali

[Install]
WantedBy=multi-user.target